# rss_downloader.service
[Unit]
Description=runs the rss_downloader

[Service]
ExecStart=/usr/bin/python /opt/bin/rss_dowloader.py

[Install]
WantedBy=multi-user.target
